# agent_deploy/package/systemd/obs-agent.service
[Unit]
Description=OBS Backup Agent
Documentation=https://obs-backup.example.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/obs-agent --service
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
User=root
Group=root
LimitNOFILE=65536
LimitNPROC=4096
Environment=OBS_AGENT_CONFIG=/etc/obs-agent/config.json

ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
